# source/doc/pycon-sk.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='application/pycon-sk' type='service' version='0'>
    <dependency name='network' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/network:default'/>
    </dependency>
    <dependency name='filesystem-local' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local:default'/>
    </dependency>
    <exec_method name='start' type='method' exec='%{daemon/bin} &amp;' timeout_seconds='60'/>
    <exec_method name='stop' type='method' exec=':kill' timeout_seconds='60'/>
    <instance name='webhooks' enabled='true'>
      <method_context working_directory='/var/tmp'>
        <method_credential group='www' user='www'/>
        <method_environment>
          <envvar name='PATH' value='/usr/local/bin:/opt/local/bin:/usr/bin:/bin'/>
          <envvar name='VIRTUAL_ENV' value='/opt/local/envs'/>
          <envvar name='BOTTLE_SERVER' value='paste'/>
          <envvar name='BOTTLE_GITHUB_SECRET' value=''/>
        </method_environment>
      </method_context>
      <property_group name='daemon' type='application'>
        <propval name='bin' type='astring' value='/var/www/cgi-bin/webhooks.py'/>
      </property_group>
    </instance>
    <stability value='Unstable'/>
  </service>
</service_bundle>
